# deploy/systemd/if_flow-archive.service
[Unit]
Description=Compress and retain old if_flow JSONL files

[Service]
Type=oneshot
EnvironmentFile=-/opt/if_flow/deploy/systemd/if_flow.env
ExecStart=/opt/if_flow/scripts/if_flow_archive.sh --base-path ${IF_FLOW_JSON_PATH:-/opt/if_flow/if_flow.jsonl} --retention-days ${IF_FLOW_RETENTION_DAYS:-7}
