# deploy/systemd/if_flow.service
[Unit]
Description=if_flow host traffic monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/opt/if_flow/deploy/systemd/if_flow.env
WorkingDirectory=/opt/if_flow
ExecStart=/opt/if_flow/scripts/run_if_flow.sh
Restart=always
RestartSec=3
KillMode=control-group
TimeoutStopSec=10
LimitMEMLOCK=infinity
LimitNOFILE=262144
StandardOutput=journal
StandardError=journal
SyslogIdentifier=if_flow

[Install]
WantedBy=multi-user.target
